# os_specific/systemd/user/budgeter_backup.timer
[Unit]
Description=Backup my budgeting information to nextcloud daily.

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
